# systemd/zino-ctx.service
[Unit]
Description=ZINO Context History Service
Documentation=file:///opt/zino/README.md
PartOf=zino.target zino-daemon.service
After=network.target zino-mem.service
Wants=zino-mem.service

[Service]
Type=simple
User=zino
Group=zino
WorkingDirectory=/opt/zino
EnvironmentFile=-/etc/zino/env
ExecStart=/opt/zino/venv/bin/python3 /opt/zino/zino-ctx.py --config /etc/zino/ZINO.toml
Restart=on-failure
RestartSec=5

# Hardening
NoNewPrivileges=yes
ProtectSystem=full
ProtectHome=yes
PrivateTmp=yes

[Install]
WantedBy=zino.target
